# systemd/cont-ai-nerd.container.in
# Podman Quadlet — cont-ai-nerd
# This is a TEMPLATE. setup.sh renders it with concrete values into
# /etc/containers/systemd/cont-ai-nerd.container
#
# Placeholders (sed-replaced by setup.sh):
#   @@AGENT_UID@@        — numeric UID of the host 'agent' user
#   @@AGENT_GID@@        — numeric GID of the host 'ai' group
#   @@VOLUME_LINES@@     — one Volume= directive per project path
#   @@PRIMARY_HOME@@     — absolute path to the primary user's home directory
#   @@CONTAINERD_CONFIG@@— path to ~/.config/cont-ai-nerd on the host
#   @@HOST@@             — listen address (default 127.0.0.1)
#   @@PORT@@             — listen port (default 3000)

[Unit]
Description=cont-ai-nerd — Sandboxed AI Agent Container
After=network-online.target

[Container]
ContainerName=cont-ai-nerd
Image=localhost/cont-ai-nerd:latest

# Run as the real host 'agent' user / 'ai' group.
# Because this is rootful podman, these map 1:1 to host UIDs.
User=@@AGENT_UID@@
Group=@@AGENT_GID@@

# ---------- Bind mounts ----------

# Project directories (read-write, SELinux relabel)
@@VOLUME_LINES@@

# cont-ai-nerd generated policy (external_directory allow-list, read-only)
Volume=@@CONTAINERD_CONFIG@@/opencode.json:/etc/cont-ai-nerd/opencode.json:ro

# cont-ai-nerd configuration (for opencode-tui wrapper)
Volume=@@CONTAINERD_CONFIG@@/config.json:/etc/cont-ai-nerd/config.json:ro

# OpenCode global configuration (read-only)
Volume=@@PRIMARY_HOME@@/.config/opencode:/home/agent/.config/opencode:ro

# OpenCode data / auth credentials (read-only)
Volume=@@PRIMARY_HOME@@/.local/share/opencode:/home/agent/.local/share/opencode:ro

# ---------- Environment ----------
# Point opencode at the cont-ai-nerd policy file (higher precedence than global config).
Environment=OPENCODE_CONFIG=/etc/cont-ai-nerd/opencode.json

# ---------- Networking ----------
# Host networking keeps things simple; opencode binds to localhost only.
Network=host

# ---------- Resource limits ----------
PodmanArgs=--memory 2g --pids-limit 100

# ---------- Entrypoint args ----------
# Run OpenCode in server mode (headless HTTP server)
Exec=serve --hostname @@HOST@@ --port @@PORT@@

[Service]
Restart=on-failure

[Install]
WantedBy=multi-user.target
