# systemd/contain-commit.timer
[Unit]
Description=Periodic cont-ai-nerd Container Commit

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
