# systemd/contain-watcher.service.in
# systemd unit — contain file-ownership watcher
# TEMPLATE: setup.sh renders this into /etc/systemd/system/
#
# Placeholders (sed-replaced by setup.sh):
#   INSTALL_DIR, PRIMARY_USER, AGENT_USER, WATCH_DIRS

[Unit]
Description=contain File Ownership Watcher
After=contain.service
Wants=contain.service

[Service]
Type=simple
ExecStart=@@INSTALL_DIR@@/contain-watcher.sh @@PRIMARY_USER@@ @@AGENT_USER@@ @@WATCH_DIRS@@
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
