# systemd/contAIn.container.in
# Podman Quadlet — contain
# This is a TEMPLATE. setup.sh renders it with concrete values into
# /etc/containers/systemd/contain.container
#
# Placeholders (sed/awk-replaced by setup.sh):
#   VOLUME_LINES, PRIMARY_HOME, CONTAINERD_CONFIG, HOST, PORT

[Unit]
Description=contain — Sandboxed AI Agent Container
After=network-online.target

[Container]
ContainerName=contAIn
Image=localhost/contain:latest

# The entrypoint runs as root to create host-path symlinks,
# then drops to the agent user via setpriv.

# ---------- Bind mounts ----------

# Project directories (read-write, SELinux relabel)
@@VOLUME_LINES@@

# contain generated policy (external_directory allow-list, read-only)
Volume=@@CONTAINERD_CONFIG@@/opencode.json:/etc/contain/opencode.json:ro

# contain configuration (for opencode-tui wrapper)
Volume=@@CONTAINERD_CONFIG@@/config.json:/etc/contain/config.json:ro

# OpenCode global configuration (read-only)
Volume=@@PRIMARY_HOME@@/.config/opencode:/home/agent/.config/opencode:ro

# OpenCode data directory (read-write) — database, auth, logs, snapshots, etc.
Volume=@@PRIMARY_HOME@@/.local/share/opencode:/home/agent/.local/share/opencode:rw

# OpenCode state directory (read-write) — model preferences, prompt history, etc.
Volume=@@PRIMARY_HOME@@/.local/state/opencode:/home/agent/.local/state/opencode:rw

# ---------- Environment ----------
# Point opencode at the contain policy file (higher precedence than global config).
Environment=OPENCODE_CONFIG=/etc/contain/opencode.json

# ---------- Networking ----------
# Host networking keeps things simple; opencode binds to localhost only.
Network=host

# ---------- Resource limits ----------
PodmanArgs=--memory 2g --pids-limit 100

# ---------- Entrypoint args ----------
# Run OpenCode in server mode (headless HTTP server)
Exec=serve --hostname @@HOST@@ --port @@PORT@@

[Service]
Restart=on-failure

[Install]
WantedBy=multi-user.target
