# systemd/contain-commit.service
# systemd unit — one-shot container image commit
# TEMPLATE: setup.sh renders @@INSTALL_DIR@@ then installs to /etc/systemd/system/

[Unit]
Description=cont-ai-nerd Container State Commit

[Service]
Type=oneshot
ExecStart=@@INSTALL_DIR@@/cont-ai-nerd-commit.sh
